# etc/systemd/system/ollama.service.d/rocm.conf
[Service]
# Expose only the RX 7800 XT (gfx1101); hide the unsupported Raphael iGPU
# (gfx1036), whose rocBLAS Tensile probe crashes during GPU discovery.
Environment="ROCR_VISIBLE_DEVICES=0"
